# deployment/supervisor/pizza.conf
[program:pizzeria]
command=/home/USER/pizzeria/bin/gunicorn --workers 3 pizzeria.wsgi:application
directory=/home/USER/django-pizzeria
user=USER
autostart=true
autorestart=true
redirect_stderr=true
